# dist/openwrt/initds/bmtools
#!/bin/sh /etc/rc.common
# "new(er)" style init script
# Look at /lib/functions/service.sh on a running system for explanations of what other SERVICE_
# options you can use, and when you might want them.

# This is a run-once script. It executes a couple of trivial benchmarks and publishes the results
# to an easy to scrape url.

START=99
STOP=0
APP=scion-control

start() {
    u=$$
    # coremark
    /usr/bin/scion-coremark > /tmp/coremark.$u.out
    if grep -q "Correct operation validated" /tmp/coremark.$u.out; then
        coremark=$(grep -s "CoreMark [0-9.]* :" /tmp/coremark.$u.out | awk '{ print $4 }' 2>/dev/null)
    fi
    if [ -z "$coremark" ]; then
        coremark=0
    fi
    rm /tmp/coremark.$u.out

    # mmbm
    mmbm=$(/usr/bin/scion-mmbm 2>/dev/null)
    if [ -z "$mmbm" ]; then
        mmbm="\"mmbm\": 0"
    fi

    # Publish
    cat <<EOF > /www/horsepower.json
{
"coremark": ${coremark},
${mmbm}
}
EOF

    # Replace the scion-router initd script with a slightly modified one...
    # (We need GOMAXPROCS set to a specific value)
    cp -f /etc/scion/router-initd.default /etc/init.d/scion-router
    chmod 755 /etc/init.d/scion-router
}

stop() {
    true
}
