# dist/openwrt/initds/router
#!/bin/sh /etc/rc.common
# "new(er)" style init script
# Look at /lib/functions/service.sh on a running system for explanations of what other SERVICE_
# options you can use, and when you might want them.

START=99
STOP=0
APP=scion-router
SERVICE_WRITE_PID=1
SERVICE_DAEMONIZE=1

init_conf() {
    # These defaults exist only if the testconfig package was installed.
    if [ -f "/etc/scion/router.toml" ]; then
        rm -f /etc/scion/router.toml.default
    elif [ -f "/etc/scion/router.toml.default" ]; then
        mv /etc/scion/router.toml.default /etc/scion/router.toml
    fi
    if [ -f "/etc/scion/topology.json" ]; then
        rm -f /etc/scion/topology.json.default
    elif [ -f "/etc/scion/topology.json.default" ]; then
        mv /etc/scion/topology.json.default /etc/scion/topology.json
    fi
    if [ -f "/etc/scion/keys/master0.key" ]; then
        rm -f /etc/scion/keys/*.key.default
    elif [ -f "/etc/scion/keys/master0.key.default" ]; then
        mv /etc/scion/keys/master0.key.default /etc/scion/keys/master0.key
        mv /etc/scion/keys/master1.key.default /etc/scion/keys/master1.key
    fi
}

start() {
    init_conf
    service_start /usr/bin/$APP --config /etc/scion/router.toml
}

stop() {
    service_stop /usr/bin/$APP
}
