# gen/supervisord.conf
[program:br1-ff00_0_110-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-ff00_0_110-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_110/br1-ff00_0_110-1.toml

[program:br1-ff00_0_110-2]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-ff00_0_110-2.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_110/br1-ff00_0_110-2.toml

[program:br1-ff00_0_110-3]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-ff00_0_110-3.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_110/br1-ff00_0_110-3.toml

[program:cs1-ff00_0_110-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs1-ff00_0_110-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/ASff00_0_110/cs1-ff00_0_110-1.toml

[program:sd1-ff00_0_110]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd1-ff00_0_110.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/ASff00_0_110/sd.toml

[group:as1-ff00_0_110]
programs = br1-ff00_0_110-1,br1-ff00_0_110-2,br1-ff00_0_110-3,cs1-ff00_0_110-1,sd1-ff00_0_110

[program:br1-ff00_0_120-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-ff00_0_120-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_120/br1-ff00_0_120-1.toml

[program:br1-ff00_0_120-2]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-ff00_0_120-2.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_120/br1-ff00_0_120-2.toml

[program:br1-ff00_0_120-3]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-ff00_0_120-3.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_120/br1-ff00_0_120-3.toml

[program:br1-ff00_0_120-4]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-ff00_0_120-4.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_120/br1-ff00_0_120-4.toml

[program:br1-ff00_0_120-5]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-ff00_0_120-5.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_120/br1-ff00_0_120-5.toml

[program:cs1-ff00_0_120-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs1-ff00_0_120-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/ASff00_0_120/cs1-ff00_0_120-1.toml

[program:sd1-ff00_0_120]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd1-ff00_0_120.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/ASff00_0_120/sd.toml

[group:as1-ff00_0_120]
programs = br1-ff00_0_120-1,br1-ff00_0_120-2,br1-ff00_0_120-3,br1-ff00_0_120-4,br1-ff00_0_120-5,cs1-ff00_0_120-1,sd1-ff00_0_120

[program:br1-ff00_0_111-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-ff00_0_111-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_111/br1-ff00_0_111-1.toml

[program:cs1-ff00_0_111-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs1-ff00_0_111-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/ASff00_0_111/cs1-ff00_0_111-1.toml

[program:sd1-ff00_0_111]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd1-ff00_0_111.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/ASff00_0_111/sd.toml

[group:as1-ff00_0_111]
programs = br1-ff00_0_111-1,cs1-ff00_0_111-1,sd1-ff00_0_111

[program:br1-ff00_0_121-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-ff00_0_121-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_121/br1-ff00_0_121-1.toml

[program:br1-ff00_0_121-2]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-ff00_0_121-2.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_121/br1-ff00_0_121-2.toml

[program:cs1-ff00_0_121-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs1-ff00_0_121-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/ASff00_0_121/cs1-ff00_0_121-1.toml

[program:sd1-ff00_0_121]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd1-ff00_0_121.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/ASff00_0_121/sd.toml

[group:as1-ff00_0_121]
programs = br1-ff00_0_121-1,br1-ff00_0_121-2,cs1-ff00_0_121-1,sd1-ff00_0_121

[program:br1-ff00_0_122-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-ff00_0_122-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_122/br1-ff00_0_122-1.toml

[program:cs1-ff00_0_122-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs1-ff00_0_122-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/ASff00_0_122/cs1-ff00_0_122-1.toml

[program:sd1-ff00_0_122]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd1-ff00_0_122.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/ASff00_0_122/sd.toml

[group:as1-ff00_0_122]
programs = br1-ff00_0_122-1,cs1-ff00_0_122-1,sd1-ff00_0_122

[program:br1-ff00_0_123-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-ff00_0_123-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_123/br1-ff00_0_123-1.toml

[program:cs1-ff00_0_123-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs1-ff00_0_123-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/ASff00_0_123/cs1-ff00_0_123-1.toml

[program:sd1-ff00_0_123]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd1-ff00_0_123.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/ASff00_0_123/sd.toml

[group:as1-ff00_0_123]
programs = br1-ff00_0_123-1,cs1-ff00_0_123-1,sd1-ff00_0_123

[program:br2-ff00_0_210-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br2-ff00_0_210-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_210/br2-ff00_0_210-1.toml

[program:br2-ff00_0_210-2]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br2-ff00_0_210-2.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_210/br2-ff00_0_210-2.toml

[program:cs2-ff00_0_210-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs2-ff00_0_210-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/ASff00_0_210/cs2-ff00_0_210-1.toml

[program:sd2-ff00_0_210]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd2-ff00_0_210.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/ASff00_0_210/sd.toml

[group:as2-ff00_0_210]
programs = br2-ff00_0_210-1,br2-ff00_0_210-2,cs2-ff00_0_210-1,sd2-ff00_0_210

[program:br2-ff00_0_211-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br2-ff00_0_211-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_211/br2-ff00_0_211-1.toml

[program:cs2-ff00_0_211-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs2-ff00_0_211-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/ASff00_0_211/cs2-ff00_0_211-1.toml

[program:sd2-ff00_0_211]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd2-ff00_0_211.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/ASff00_0_211/sd.toml

[group:as2-ff00_0_211]
programs = br2-ff00_0_211-1,cs2-ff00_0_211-1,sd2-ff00_0_211

[program:br3-ff00_0_310-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br3-ff00_0_310-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_310/br3-ff00_0_310-1.toml

[program:br3-ff00_0_310-2]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br3-ff00_0_310-2.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_310/br3-ff00_0_310-2.toml

[program:br3-ff00_0_310-3]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br3-ff00_0_310-3.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_310/br3-ff00_0_310-3.toml

[program:cs3-ff00_0_310-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs3-ff00_0_310-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/ASff00_0_310/cs3-ff00_0_310-1.toml

[program:sd3-ff00_0_310]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd3-ff00_0_310.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/ASff00_0_310/sd.toml

[group:as3-ff00_0_310]
programs = br3-ff00_0_310-1,br3-ff00_0_310-2,br3-ff00_0_310-3,cs3-ff00_0_310-1,sd3-ff00_0_310

[program:br3-ff00_0_311-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br3-ff00_0_311-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_311/br3-ff00_0_311-1.toml

[program:cs3-ff00_0_311-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs3-ff00_0_311-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/ASff00_0_311/cs3-ff00_0_311-1.toml

[program:sd3-ff00_0_311]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd3-ff00_0_311.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/ASff00_0_311/sd.toml

[group:as3-ff00_0_311]
programs = br3-ff00_0_311-1,cs3-ff00_0_311-1,sd3-ff00_0_311

[program:br4-ff00_0_410-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br4-ff00_0_410-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_410/br4-ff00_0_410-1.toml

[program:br4-ff00_0_410-2]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br4-ff00_0_410-2.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_410/br4-ff00_0_410-2.toml

[program:cs4-ff00_0_410-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs4-ff00_0_410-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/ASff00_0_410/cs4-ff00_0_410-1.toml

[program:sd4-ff00_0_410]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd4-ff00_0_410.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/ASff00_0_410/sd.toml

[group:as4-ff00_0_410]
programs = br4-ff00_0_410-1,br4-ff00_0_410-2,cs4-ff00_0_410-1,sd4-ff00_0_410

[program:br4-ff00_0_411-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br4-ff00_0_411-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_411/br4-ff00_0_411-1.toml

[program:cs4-ff00_0_411-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs4-ff00_0_411-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/ASff00_0_411/cs4-ff00_0_411-1.toml

[program:sd4-ff00_0_411]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd4-ff00_0_411.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/ASff00_0_411/sd.toml

[group:as4-ff00_0_411]
programs = br4-ff00_0_411-1,cs4-ff00_0_411-1,sd4-ff00_0_411

[program:br5-ff00_0_510-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br5-ff00_0_510-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/ASff00_0_510/br5-ff00_0_510-1.toml

[program:cs5-ff00_0_510-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs5-ff00_0_510-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/ASff00_0_510/cs5-ff00_0_510-1.toml

[program:sd5-ff00_0_510]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd5-ff00_0_510.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/ASff00_0_510/sd.toml

[group:as5-ff00_0_510]
programs = br5-ff00_0_510-1,cs5-ff00_0_510-1,sd5-ff00_0_510

[program:dispatcher]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/dispatcher.log
redirect_stderr = True
startretries = 0
startsecs = 1
priority = 50
command = bin/dispatcher --config gen/dispatcher/disp.toml

