# dist/systemd/scion-ip-gateway.service
[Unit]
Description=SCION IP Gateway
Documentation=https://docs.scion.org
After=network-online.target scion-daemon.service
Wants=scion-daemon.service

[Service]
Type=simple
User=scion
Group=scion
ExecStart=/usr/bin/scion-ip-gateway --config /etc/scion/sig.toml
Restart=on-failure

[Install]
WantedBy=multi-user.target
