# dist/systemd/scion-control@.service
[Unit]
Description=SCION Control Service
Documentation=https://docs.scion.org
After=network-online.target scion-dispatcher.service
Wants=scion-dispatcher.service
StartLimitBurst=1
StartLimitInterval=1s

[Service]
Type=simple
User=scion
Group=scion
ExecStart=/usr/bin/scion-control --config /etc/scion/%i.toml
Restart=on-failure

[Install]
WantedBy=multi-user.target
