# dist/systemd/scion-router@.service
[Unit]
Description=SCION Router
Documentation=https://docs.scion.org
After=network-online.target
StartLimitBurst=1
StartLimitInterval=1s

[Service]
Type=simple
User=scion
Group=scion
ExecStart=/usr/bin/scion-router --config /etc/scion/%i.toml
Restart=on-failure
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_BPF

[Install]
WantedBy=multi-user.target
